# roles/rdo-infra/bot/templates/rdobot.service.j2
[Unit]
Description=rdobot
After=syslog.target network.target

[Service]
ExecStart={{ rdobot_systemd_wrapper }}
ExecStop=/bin/kill -s INT $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target